# schedule.cron
30 * * * * . /app/project_env.sh; cd /app/ && /app/exec.sh >> /var/log/cron.log 2>&1
# remember to end this file with an empty new line